# deploy/utrecht-bus-web.service
[Unit]
Description=Utrecht bus monitor - webserver (gunicorn)
After=network.target utrecht-bus-collector.service

[Service]
Type=simple
User=utrechtbus
WorkingDirectory=/opt/utrecht-bus-monitor
EnvironmentFile=-/opt/utrecht-bus-monitor/.env
ExecStart=/opt/utrecht-bus-monitor/venv/bin/gunicorn --preload -w 2 --timeout 120 -b 127.0.0.1:5151 app.wsgi:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
